# buildroot/package/glowbarn-sensors/glowbarn-sensors.service
[Unit]
Description=GlowBarn Sensors Service
Documentation=https://github.com/bad-antics/glowbarn-os
After=network.target local-fs.target
Requires=glowbarn.service
After=glowbarn.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 -m glowbarn.sensors
WorkingDirectory=/opt/glowbarn
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONPATH=/opt/glowbarn/lib
Environment=GLOWBARN_CONFIG=/etc/glowbarn/config.yaml

[Install]
WantedBy=multi-user.target
